# images/cron/crontab
* * * * * /usr/local/bin/php /application/artisan schedule:run >> /application/storage/logs/cron.log
